# prosody.conf
[supervisord]
nodaemon=true

[program:prosody]
command=prosody
user=prosody
